# deploy/integrations/vps_legacy/portable_bundle/legacy/server/watchdog/toy-watchdog.service
[Unit]
Description=PrivateTunnel toy UDP/TUN gateway watchdog
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/privatetunnel/toy-gateway
ExecStart=/usr/bin/python3 /opt/privatetunnel/toy-gateway/toy_tun_gateway.py
Restart=on-failure
RestartSec=3
# Capture stdout/stderr to journal for post-mortem analysis
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
